# common/etc/systemd/system/ldm.service.d/override.conf
[Service]
EnvironmentFile=
ExecStart=
ExecStart=/usr/bin/ldm -u Link -p /media -c /home/Link/.bin/ldm.sh
