# home/Link/.config/systemd/user/appskey-is-F13.service
[Unit]
Description=Appskey is F13

[Service]
Type=oneshot
ExecStart=/usr/bin/loadkeys /home/Link/.bin/F13.map

[Install]
WantedBy=multi-user.target
