# home/Link/.config/systemd/user/clipmenud.service
[Unit]
Description=Clipmenud

[Service]
Type=simple
ExecStart=/usr/bin/clipmenud
Restart=always

[Install]
WantedBy=default.target
