# home/Link/.config/systemd/user/dwmstatus.timer
[Unit]
Description=dwm statusbar

[Timer]
OnUnitActiveSec=10s
Unit=dwmstatus.service

[Install]
WantedBy=default.target
