# home/Link/.config/systemd/user/feh-wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnUnitActiveSec=180s
Unit=feh-wallpaper.service

[Install]
WantedBy=multi-user.target
