# home/Link/.config/systemd/user/feh-wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnUnitActiveSec=720
Unit=feh-wallpaper.service

[Install]
WantedBy=default.target
