# common/home/user/.config/systemd/user/compton.service
[Unit]
Description=compton

[Service]
Type=simple
ExecStart=/usr/bin/compton --config /tmp/compton.conf
Restart=always
RestartSec=30

[Install]
WantedBy=x11.target
Requisite=x11.target
