# common/home/user/.config/systemd/user/rawdog.service
[Unit]
Description=Update rawdog feeds

[Service]
Type=oneshot
ExecStart=/usr/bin/rawdog -uw -d %h/.config/rawdog
