# common/home/user/.config/systemd/user/rawdog.timer
[Unit]
Description=Update rawdog feeds

[Timer]
OnUnitActiveSec=3600
Unit=rawdog.service
