# common/home/user/.config/systemd/user/sleep.target
[Unit]
Description=Sleep Timer Target
StopWhenUnneeded=yes
