# common/home/user/.config/systemd/user/ssh-notify.service
[Unit]
Description=SSH irc notifications

[Service]
Type=simple
ExecStart=/home/Link/.bin/ssh-notify
Restart=always
RestartSec=10

[Install]
WantedBy=x11.target
Requisite=x11.target
