# common/home/user/.config/systemd/user/suspend.target
[Unit]
Description=Suspend Timer Target
StopWhenUnneeded=yes
