# common/home/user/.config/systemd/user/x11.target
[Unit]
Description=X is running
StopWhenUnneeded=yes
