# GIR/etc/systemd/system/appskey-is-tab.service
[Unit]
Description=Appskey is tab

[Service]
Type=oneshot
ExecStart=/usr/bin/setkeycodes e05d 15

[Install]
WantedBy=multi-user.target
