# systemd/.config/systemd/user/sunshine.service.d/override.conf
[Service]
Environment=DISPLAY=:0.0
ExecStart=
ExecStart=/usr/bin/sunshine %E/sunshine/sunshine.conf

[Install]
WantedBy=default.target
