# systemd/.config/systemd/user/redshift.service.d/override.conf
[Unit]
Requisite=x11.target
After=x11.target

[Install]
WantedBy=
WantedBy=x11.target
