# systemd/.config/systemd/user/feh-wallpaper.timer
[Unit]
Description=Random wallpaper with feh
Requisite=x11.target

[Timer]
OnStartupSec=0
OnUnitActiveSec=720

[Install]
WantedBy=x11.target
