# systemd/.config/systemd/user/rawdog.timer
[Unit]
Description=Update rawdog feeds

[Timer]
OnUnitActiveSec=3600

[Install]
WantedBy=default.target
