# systemd/.config/systemd/user/sunshine.service
[Unit]
Description=Sunshine Gamestream Server for Moonlight

[Service]
ExecStart=/usr/bin/sunshine

[Install]
WantedBy=graphical-session.target
