# systemd/.config/systemd/user/x11.target
[Unit]
Description=X11 target
StopWhenUnneeded=yes
