# common/home/user/.config/systemd/user/dunst.service.d/override.conf
[Unit]
Requisite=x11.target

[Install]
WantedBy=x11.target
