# Talim/etc/systemd/system/talim-extra-keys2.service
[Unit]
Description=Assign Talim's Extra Keys

[Service]
Type=oneshot
ExecStart=/usr/bin/setkeycodes e007 108

[Install]
WantedBy=multi-user.target
